# packaging/openwrt/luci/root/etc/init.d/ruhop
#!/bin/sh /etc/rc.common
# Copyright 2024 ruhop contributors
# Licensed under the Apache License, Version 2.0
#
# Ruhop VPN init script for OpenWRT with LuCI integration

USE_PROCD=1
START=99
STOP=01

PROG=/usr/bin/ruhop
CONFIG_GEN=/usr/lib/ruhop/gen-config.sh

validate_section_ruhop() {
    uci_load_validate ruhop ruhop "$1" "$2" \
        'enabled:bool:0' \
        'mode:string:client' \
        'config:string:/etc/ruhop/ruhop.toml'
}

start_service() {
    local enabled mode config

    validate_section_ruhop main start_instance
}

start_instance() {
    local enabled="$1"
    local mode="$2"
    local config="$3"

    [ "$enabled" = "0" ] && return 0

    # Generate TOML config from UCI settings
    if [ -x "$CONFIG_GEN" ]; then
        "$CONFIG_GEN" "$config"
    fi

    # Ensure TUN module is loaded
    [ ! -d /dev/net ] && mkdir -p /dev/net
    [ ! -c /dev/net/tun ] && mknod /dev/net/tun c 10 200 2>/dev/null
    modprobe tun 2>/dev/null

    procd_open_instance ruhop
    procd_set_param command "$PROG" "$mode" --config "$config"

    # Respawn settings: respawn after crash
    # threshold=3600 (seconds) - period to track crashes
    # timeout=5 (seconds) - delay before respawn
    # retry=5 - max respawn attempts within threshold
    procd_set_param respawn 3600 5 5

    # File watch - reload on config changes
    procd_set_param file "$config"
    procd_set_param file /etc/config/ruhop

    # Resource limits
    procd_set_param limits nofile=65535

    # Logging
    procd_set_param stdout 1
    procd_set_param stderr 1

    procd_close_instance
}

stop_service() {
    # procd handles stopping automatically
    :
}

reload_service() {
    # Regenerate config and restart
    stop
    start
}

service_triggers() {
    procd_add_reload_trigger "ruhop"
    procd_add_validation validate_section_ruhop
}
